# misc/supervisor.ini
[program:mdb-links]
process_name = mdb-links
command = /sites/mdb-links/mdb-links server
directory = /sites/mdb-links/
user = archive
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /sites/mdb-links/logs/links.log